# deploy/freeq-server.service
[Unit]
Description=FreeQ IRC Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chad
Group=chad
WorkingDirectory=/home/chad/src/freeq
ExecStart=/home/chad/src/freeq/target/release/freeq-server \
    --listen-addr 0.0.0.0:6667 \
    --tls-listen-addr 0.0.0.0:6697 \
    --tls-cert /etc/letsencrypt/live/tech.blueyard.com/fullchain.pem \
    --tls-key /etc/letsencrypt/live/tech.blueyard.com/privkey.pem \
    --db-path ./irc.db \
    --web-addr 0.0.0.0:8080 \
    --web-static-dir /home/chad/src/freeq/freeq-app/dist \
    --server-name irc.freeq.at \
    --motd "Welcome to freeq — IRC with AT Protocol identity. https://freeq.at" \
    --iroh
EnvironmentFile=/home/chad/src/freeq/.env.secrets
Restart=on-failure
RestartSec=5

# Allow binding to low ports (6667, 6697, 8080)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/chad/src/freeq
ReadOnlyPaths=/etc/letsencrypt

[Install]
WantedBy=multi-user.target
